# systemd/brandyfication-mcp.service
[Unit]
Description=BRANDYFICATION File Hosting MCP Server
Documentation=https://github.com/your-repo/brandyfication
After=network.target

[Service]
Type=simple
User=brandyfication
Group=brandyfication
WorkingDirectory=/opt/brandyfication
Environment=NODE_ENV=production
Environment=STORAGE_DIR=/var/lib/brandyfication/storage
ExecStart=/usr/bin/node /opt/brandyfication/dist/index.js
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=brandyfication-mcp

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/brandyfication

[Install]
WantedBy=multi-user.target
