# systemd/brandyfication-agent.service
[Unit]
Description=BRANDYFICATION Agent Frontend
Documentation=https://github.com/your-repo/brandyfication
After=network.target brandyfication-mcp.service
Wants=brandyfication-mcp.service

[Service]
Type=simple
User=brandyfication
Group=brandyfication
WorkingDirectory=/opt/brandyfication/agent
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=STORAGE_DIR=/var/lib/brandyfication/storage
ExecStart=/usr/bin/node /opt/brandyfication/agent/dist/agent.js
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=brandyfication-agent

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/brandyfication

[Install]
WantedBy=multi-user.target
